# conf/AltServer.target
[Unit]
Description=__APP__
Wants=avahi-daemon.service usbmuxd.service AltServer.service